# ubuntu-laravel-5x-wkhtmltopdf/configs/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisor/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)

[unix_http_server]
file=/var/run/supervisor/supervisord.sock   ; (the path to the socket file)

[program:nginx]
command=nginx ;
stdout_logfile=/dev/stdout ;
stdout_logfile_maxbytes=0 ;
stderr_logfile=/dev/stderr ;
stderr_logfile_maxbytes=0 ;
autostart=true ;
autorestart=true ;

[program:php-fpm]
command=php-fpm7.2 ;
stdout_logfile=/dev/stdout ;
stdout_logfile_maxbytes=0 ;
stderr_logfile=/dev/stderr ;
stderr_logfile_maxbytes=0 ;
autostart=true ;
autorestart=true ;